# oat-evaluation/light-controller-cb/telldusd.service
[Unit]
Description=Telldus daemon

[Service]
Type=simple
ExecStart=/usr/sbin/telldusd --nodaemon
Restart=always

[Install]
WantedBy=multi-user.target
